# omarchy/.config/systemd/user/tmux-auto-save.service
[Unit]
Description=Tmux session auto-save via resurrect

[Service]
Type=oneshot
ExecStart=%h/.local/bin/tmux-auto-save
